# deploy/openclaw-personal.service
[Unit]
Description=OpenClaw Personal Bot
After=network.target
Wants=network-online.target

[Service]
Type=exec
User=openclaw
Group=openclaw
WorkingDirectory=/opt/openclaw
ExecStart=/usr/local/bin/doppler run --config personal --project openclaw -- /opt/openclaw/.venv/bin/python -m openclaw
Restart=always
RestartSec=10
MemoryMax=2G
CPUQuota=150%

[Install]
WantedBy=multi-user.target
