# deploy/openclaw-dev.service
[Unit]
Description=OpenClaw Dev Bot
After=network.target
Wants=network-online.target

[Service]
Type=exec
User=openclaw
Group=openclaw
WorkingDirectory=/opt/openclaw-dev
ExecStart=/usr/local/bin/doppler run --config dev --project openclaw -- /opt/openclaw-dev/.venv/bin/python -m openclaw
Restart=always
RestartSec=10
MemoryMax=1G
CPUQuota=100%

[Install]
WantedBy=multi-user.target
